# supervisor.ini
[program:tictactoe]
command=/home/h0ffmann/bin/tictactoe.py
startsecs=60

